# service/seek_camera.service
[Unit]
Description= open seek camera

[Service]
ExecStart= /var/thermal/start_seek.sh 
Restart=always
RestartSec=3

[Install]
WantedBy=default.target